# light.service
[Unit]
Description=Lightd

[Service]
Type=simple
ExecStart=/usr/bin/lightd

[Install]
WantedBy=multi-user.target
