# light.service
[Unit]
Description=Lightd
[Service]
ExecStart=/usr/bin/lightd

[Install]
WantedBy=default.target
